# angellevel-plasma/tools/systemd/angellevel-update-notifier.service
[Unit]
Description=AngelLevel software-update notifier
Documentation=man:angellevel-update-notifier

[Service]
Type=oneshot
ExecStart=%h/.local/bin/angellevel-update-notifier --quiet
